# roles/postgresql/templates/etc/profile.d/postgres.sh.j2
# {{ ansible_managed }}

export PGDATA={{ postgres_data_dir }}
export PATH=$PATH:{{ postgres_bin_dir }}
